# script/app.supervisor.conf
[program:${PROGRAM}]
command=${CMD}
directory=${DIR}
autostart=true
autorestart=true
startretries=3
stderr_logfile=${STDERR}
stdout_logfile=${STDOUT}